# systemd/cellaserv2.service
[Unit]
Description=Cellaserv2 RPC
PartOf=evolutek.target

[Service]
Restart=always
User=evolutek
ExecStart=/usr/bin/cellaserv2
