# systemd/evolutek.target
[Unit]
Description=Evolutek
Requires=cellaserv.service
